# docker/images/ezstream/config/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:ezstream]
command=ezstream -c /etc/ezstream.xml
stdout_logfile=/var/log/ezsteam.log
autorestart=true

